# host/crontab
0 7 1 * * /home/haris/website/.tem/path/reissue-cert > /dev/null
